# res/rc.local
#!/bin/bash
/data/bridge/tailscale/tailscaled --state=/var/lib/tailscale/tailscaled.state --socket=/run/tailscale/tailscaled.sock --port=41641 &
 sleep 20 && /data/bridge/startup.sh &